# systemd/kimsufi-notifier.timer
[Unit]
Description=Run kimsufi-notifier minutely

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true

[Install]
WantedBy=timers.target
